# system_files/ad/usr/lib/systemd/system/sddm.service.d/tholdyos.conf
[Unit]
After=flatpak-preinstall.service
After=tholdyos-domain-setup.service
